# luci-app-sms-tool-js/root/etc/uci-defaults/off_sms.sh
#!/bin/sh
# Copyright 2023 Rafał Wabik (IceG) - From eko.one.pl forum
# Licensed to the GNU General Public License v3.0.

/etc/init.d/my_new_sms stop
/etc/init.d/my_new_sms disable

 
exit 0
